# pi-ui.service
[Unit]
Description=pi-lot-ui - web UI for the pi coding agent

[Service]
WorkingDirectory=__PI_UI_DIR__
ExecStart=/usr/bin/node __PI_UI_DIR__/server.mjs
# always: a stray SIGTERM/clean exit (e.g. an agent running pkill) must not
# leave the UI down
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
